# deploy/simple-worker.service
# =============================================================================
# systemd service for AitherHub Queue Worker
# =============================================================================
# Install:
#   sudo cp deploy/simple-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable simple-worker.service
#   sudo systemctl start simple-worker.service
#
# Logs:
#   journalctl -u simple-worker.service -f
# =============================================================================

[Unit]
Description=AitherHub Queue Worker
After=network.target
Wants=worker-health.service

[Service]
Type=simple
User=azureuser
Group=azureuser
WorkingDirectory=/var/www/aitherhub/worker/batch
Environment=PYTHONPATH=/var/www/aitherhub:/var/www/aitherhub/worker/batch
EnvironmentFile=/var/www/aitherhub/.env

# Use the NEW independent entrypoint
ExecStart=/opt/aitherhub/.venv/bin/python -m worker.entrypoints.queue_worker

# Graceful shutdown
KillSignal=SIGTERM
KillMode=process
TimeoutStopSec=300

# Auto-restart on failure
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitherhub-worker

[Install]
WantedBy=multi-user.target
