# deploy/worker-health.service
# =============================================================================
# systemd service for AitherHub Worker Health Check
# =============================================================================
# Runs alongside the queue worker, providing health check endpoints.
#
# Install:
#   sudo cp deploy/worker-health.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable worker-health.service
#   sudo systemctl start worker-health.service
#
# Test:
#   curl http://localhost:8081/health
# =============================================================================

[Unit]
Description=AitherHub Worker Health Check
After=network.target

[Service]
Type=simple
User=azureuser
Group=azureuser
# v6: ALL paths unified to /opt/aitherhub (was /var/www/aitherhub — stale copy)
WorkingDirectory=/opt/aitherhub
Environment=PYTHONPATH=/opt/aitherhub:/opt/aitherhub/worker/batch
EnvironmentFile=/opt/aitherhub/.env

ExecStart=/opt/aitherhub/.venv/bin/python -m worker.entrypoints.health_check

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitherhub-worker-health

[Install]
WantedBy=multi-user.target
